# crontab
0 5 * *	0 /etc/haproxy/update_cert.sh >> /var/log/haproxy/certbot.log 2>&1